# Intel-bay-trail-freeze-workaround/DEBIAN/postinst
#!/bin/sh
set -e

     
    systemctl daemon-reload
    systemctl enable c6offc7on.service

exit 0
