# Intel-bay-trail-freeze-workaround/etc/systemd/system/c6offc7on.service
[Unit]
Description=Disable C6 and enable C7

[Service]
ExecStart=/usr/bin/c6offc7on.sh

[Install]
WantedBy=multi-user.target
